# roles/xray/templates/xray.service.j2
[Unit]
Description=Xray Service
Documentation=https://github.com/XTLS/Xray-core
After=network.target nss-lookup.target

[Service]
User={{ xray_user }}
Group={{ xray_group }}
Type=simple
ExecStart={{ xray_install_dir }}/xray run -config {{ xray_config_dir }}/config.json
Restart=on-failure
RestartSec=5s
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
